# restic/crontab
15 5 * * *      restic    $HOME/bin/restic-backup.sh $HOME/jobs.json 2>&1 | logger -t "restic-backup"